# systemd/user/bedtime.service

[Unit]
Description=22:00 bedtime reminder

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/bedtime.sh
