# systemd/user/bedtime.timer

[Unit]
Description=Trigger bedtime reminder at 22:00 daily

[Timer]
OnCalendar=*-*-* 22:00:00

[Install]
WantedBy=timers.target
